# roles/pdns_recursor/templates/opennic-zones.timer.j2
[Unit]
Description=Update Opennic zones every hour

[Timer]
OnCalendar=hourly
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target
