# roles/pdns_recursor/templates/stevenblack-hosts.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Update StevenBlack hosts every week

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
